# deploy/systemd/basin-kiosk.service
[Unit]
Description=Basin Kiosk Display
After=graphical.target basin-api.service
Wants=graphical.target

[Service]
Type=simple
User=basin
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/basin/.Xauthority
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/xinit /usr/bin/chromium-browser \
    --kiosk \
    --no-first-run \
    --no-sandbox \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --noerrdialogs \
    http://localhost/status \
    -- :0 -nolisten tcp
Restart=on-failure
RestartSec=15

[Install]
WantedBy=graphical.target
