# deploy/systemd/basin-beat.service
[Unit]
Description=Basin Celery Beat Scheduler
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=basin
WorkingDirectory=/opt/basin/repo/backend
EnvironmentFile=/opt/basin/repo/backend/.env
ExecStart=/opt/basin/venv/bin/celery -A app.tasks.celery_app beat --loglevel=info --schedule=/opt/basin/data/celerybeat-schedule
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
